# scripts/gateway-watchdog.service
[Unit]
Description=Gateway Active Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gateway-watchdog.sh
User=root
Group=root
