# scripts/sync-offsite-backups.timer
[Unit]
Description=Run Offsite Backup Sync Daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
